# cookbooks/mailhog/templates/mailhog.service.erb
[Unit]
Description=MailHog service

[Service]
ExecStart=/usr/local/bin/mailhog \
  -api-bind-addr 0.0.0.0:<%= @settings["port"] %> \
  -ui-bind-addr 0.0.0.0:<%= @settings["port"] %> \
  -smtp-bind-addr 0.0.0.0:<%= @settings["smtp_port"] %>

[Install]
WantedBy=multi-user.target
